# subsystem/statTools/cron.d/cricket.in
## info: Runs the cricket collector, and a nightly makecricketconfig run.
*/5 * * * * @CRICKETDIR@/collect-subtrees normal 
12 5 * * * @bindir@/makecricketconfig.pl
0 5 * * * @bindir@/cleanrrds.pl -d -p @localstatedir@/cricket-data/
