# python/nav/etc/cron.d/dbclean
## info: Clean database with navclean
*/5 * * * * navclean -f -q --netbox --websessions
