# python/nav/etc/cron.d/emailreports
## info: Send email reports for device and link availability
0 5 1 * * emailreports month
0 5 * * mon emailreports week
0 5 * * * emailreports day
