# python/nav/etc/cron.d/maintengine
## info: Regularly check the maintenance-queue and send events to eventq
*/5 * * * * maintengine.py
