# python/nav/etc/cron.d/navdhcpstats
## info: Regularly fetch DHCP stats from endpoints in dhcpstats.conf to the Carbon backend
*/5 * * * * navdhcpstats
