# python/nav/etc/cron.d/navoui
## info: Update the database with OUIs and their respective vendors
0 10 * * * navoui
