# python/nav/etc/cron.d/netbiostracker
## info: Regularly fetch netbiosnames from active computers
*/15 * * * * netbiostracker.py
