# python/nav/etc/cron.d/sortedstats_cacher
## info: Regularly update cache for sortedstats
*/15 * * * * sortedstats_cacher hour
0 * * * * sortedstats_cacher day
0 */6 * * * sortedstats_cacher week
0 0 * * * sortedstats_cacher month
