# python/nav/etc/cron.d/thresholdmon
## info: Monitors your Graphite metrics for exceeded thresholds

*/5 * * * * thresholdmon
