# subsystem/dbTools/backup.cron
## info: Runs a nightly backup of the PostgreSQL database
# Backups are also rotated
58 4 * * * $BINDIR/pg_backup.pl
