# etc/cron.d/activeip
## info: Collect active ip-adresses for all prefixes
*/30 * * * * collect_active_ip.py
