# etc/cron.d/activeip.in
## info: Collect active ip-adresses for all prefixes
*/30 * * * * @bindir@/collect_active_ip.py
