# etc/cron.d/dbclean.in
## info: Clean database with navclean
*/5 * * * * @bindir@/navclean -f -q --netbox --websessions
