# etc/cron.d/emailreports.in
## info: Send email reports for device and link availability
0 5 1 * * @bindir@/report_sender.py month
0 5 * * mon @bindir@/report_sender.py week
0 5 * * * @bindir@/report_sender.py day
