# etc/cron.d/emailreports
## info: Send email reports for device and link availability
0 5 1 * * emailreports.py month
0 5 * * mon emailreports.py week
0 5 * * * emailreports.py day
