# etc/cron.d/emailreports.in
## info: Send email reports for device and link availability
0 5 1 * * @bindir@/emailreports.py month
0 5 * * mon @bindir@/emailreports.py week
0 5 * * * @bindir@/emailreports.py day
