# etc/cron.d/mactrace
## info: Checks NAV's cam log for watched MAC addresses
11,26,41,56 * * * *      macwatch.py
