# etc/cron.d/navstats
## navstats: Collect statistics defined in navstats.conf
*/5 * * * * navstats
