# etc/cron.d/navstats.in
## navstats: Collect statistics defined in navstats.conf
*/5 * * * * @bindir@/navstats
