# etc/cron.d/netbiostracker
## info: Regularly fetch netbiosnames from active computers
*/15 * * * * @bindir@/netbiostracker.py
