# etc/cron.d/networkDiscovery.in
## info: Detects the topology of your network

35 * * * *              @bindir@/navtopology
38 * * * *              @bindir@/networkDiscovery.sh vlan
