# etc/cron.d/psuwatch
## info: Monitors the state of redundant PSUs and fans

5 * * * *	powersupplywatch.py
