# etc/cron.d/stats.in
## stats: Collect statistics defined in stats.conf
*/10 * * * * @bindir@/stats.py
