# etc/cron.d/thresholdMon.in
## info: Monitors your RRD data for exceeded threshholds

*/5 * * * * @bindir@/thresholdMon.py >> @localstatedir@/log/thresholdMon.log
