# etc/cron.d/topology.in
## info: Detects the topology of your network.

35 * * * *              @bindir@/navtopology --l2 --vlan

# If you want to run the old topology detector, uncomment the following lines:

# 35 * * * *              @bindir@/networkDiscovery.sh topology
# 38 * * * *              @bindir@/networkDiscovery.sh vlan
