# subsystem/smsd/safe_smsd.cron
## info: Restarts the SMS Daemon when it has gone away
*/5 * * * * $BINDIR/safe_smsd.pl
