# subsystem/thresholdMon/cron.d/thresholdMon
## info: Monitors your RRD data for exceeded threshholds

*/5 * * * * $BINDIR/thresholdMon.py >> $LOCALSTATEDIR/log/thresholdMon.log
